# etc/sv/uwsgi/run
#!/bin/sh
exec 2>&1
exec uwsgi --emperor /etc/uwsgi/apps-enabled
